# deployment/abmg-daemon.service
[Unit]
Description=ABMG - AudioBookBay Magnet Grabber Daemon
Documentation=https://github.com/StarlightDaemon/ABMG
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/abmg
Environment="PYTHONPATH=/opt/abmg"
ExecStart=/opt/abmg/venv/bin/python -m src.abmg daemon --interval 60
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
