# deployment/abmg-rss.service
[Unit]
Description=ABMG RSS Feed Server
Documentation=https://github.com/StarlightDaemon/ABMG
After=network.target abmg-daemon.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/abmg
Environment="PYTHONPATH=/opt/abmg"
Environment="FLASK_ENV=production"
Environment="MAGNETS_FILE=/opt/abmg/magnets.txt"
ExecStart=/opt/abmg/venv/bin/python -m src.rss_server
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
