# deployment/bind.service
[Unit]
Description=BIND - Book Indexing Network Daemon
Documentation=https://github.com/StarlightDaemon/BIND
After=network.target

[Service]
Type=simple
User=bind
Group=bind
WorkingDirectory=/opt/bind
Environment="PYTHONPATH=/opt/bind"
EnvironmentFile=-/opt/bind/config.env
ExecStart=/opt/bind/venv/bin/python -m src.bind daemon
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ReadWritePaths=/opt/bind/data /opt/bind/logs /opt/bind/magnets
CapabilityBoundingSet=

[Install]
WantedBy=multi-user.target
